# dist/debian.conffiles
/etc/cron.d/rapidblock
/etc/default/rapidblock
